# packaging/templates/systemd/keinfs-kst@.service.in
[Unit]
Description=KeInFS Storage Target %i
After=network-online.target local-fs.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
PermissionsStartOnly=true
EnvironmentFile={{SYSCONFDIR}}/keinfs/kst/%i.env
ExecStartPre=/usr/bin/install -d {{RUNSTATEDIR}}/keinfs/kst {{RUNSTATEDIR}}/keinfs/kix
ExecStart={{LIBEXECDIR}}/kst-runner
Restart=on-failure
RestartSec=2
LimitNOFILE=1048576
WorkingDirectory={{LOCALSTATEDIR}}/lib

[Install]
WantedBy=multi-user.target
